# contrib/init/metrixd.conf
description "Metrix Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env METRIXD_BIN="/usr/bin/metrixd"
env METRIXD_USER="Metrix"
env METRIXD_GROUP="Metrix"
env METRIXD_PIDDIR="/var/run/metrixd"
# upstart can't handle variables constructed with other variables
env METRIXD_PIDFILE="/var/run/metrixd/metrixd.pid"
env METRIXD_CONFIGFILE="/etc/Metrix/Metrix.conf"
env METRIXD_DATADIR="/var/lib/metrixd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # metrixd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$METRIXD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run metrixd."
        echo "The setting must appear in $METRIXD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$METRIXD_PIDDIR"
    chmod 0755 "$METRIXD_PIDDIR"
    chown $METRIXD_USER:$METRIXD_GROUP "$METRIXD_PIDDIR"
    chown $METRIXD_USER:$METRIXD_GROUP "$METRIXD_CONFIGFILE"
    chmod 0660 "$METRIXD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$METRIXD_PIDFILE" \
    --chuid $METRIXD_USER:$METRIXD_GROUP \
    --exec "$METRIXD_BIN" \
    -- \
    -pid="$METRIXD_PIDFILE" \
    -conf="$METRIXD_CONFIGFILE" \
    -datadir="$METRIXD_DATADIR" \
    -disablewallet \
    -daemon